# meet.service
[Unit]
Description=Meet Video Meeting App
After=network.target

[Service]
Type=simple
WorkingDirectory=/path/to/your/meet
ExecStart=/usr/local/bin/bun run start
Restart=always
RestartSec=5
EnvironmentFile=/path/to/your/meet/.env
User=www-data

[Install]
WantedBy=multi-user.target 